# payments_demo/conf/etc/supervisor/conf.d/gunicorn.conf
[program:mobilepay-ui]
command=/usr/bin/gunicorn app:app -b localhost:5000
directory=/code
